# HorseIsleServer/HISPd/Resources/DEBIAN/postinst
#!/bin/sh

useradd horseisle
chown -R horseisle /etc/hisp
systemctl enable HISP.service
systemctl start HISP

echo "Make sure to edit /etc/hisp/server.properties!"
