# HorseIsleServer/HISPd/Resources/DEBIAN/postrm
#!/bin/sh
systemctl daemon-reload
systemctl stop HISP
systemctl disable HISP
userdel horseisle
