# HorseIsleServer/HISPd/Resources/DEBIAN/prerm
#!/bin/sh
systemctl stop HISP
systemctl disable HISP
userdel horseisle
